# deploy/toolbasecamp-api.service
[Unit]
Description=Tool Basecamp API (FastAPI)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/toolbasecamp-api
Environment=HOST=127.0.0.1
Environment=PORT=8001
EnvironmentFile=-/etc/toolbasecamp-api.env
ExecStart=/opt/toolbasecamp-api/venv/bin/python -B run.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
